# data/supervisor/breadcount.conf
; /etc/supervisor/conf.d/breadcount.conf

[program:breadcount-ros2]
; Use Absolute Path for the launch file to prevent "file not found" errors
command=/bin/bash -c "source /opt/tros/humble/setup.bash && ros2 launch /home/sunrise/BreadCounting/src/ros2/Ros2PipelineLauncher.py"
directory=/home/sunrise/BreadCounting/
user=sunrise
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/supervisor/ros2_err.log
stdout_logfile=/var/log/supervisor/ros2_out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
environment=HOME="/home/sunrise"

[program:breadcount-main]
; This one looked fine, but ensure main.py is strictly in BreadCounting root
command=/bin/bash -c "source /opt/tros/humble/setup.bash && /usr/bin/python main.py"
directory=/home/sunrise/BreadCounting
user=sunrise
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/main_out.log
stderr_logfile=/var/log/supervisor/main_err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5


[program:breadcount-uvicorn]
; --- FIX 1: Set directory to the PROJECT ROOT ---
directory=/home/sunrise/BreadCounting
; --- FIX 2: Run uvicorn as a module from the root ---
; Python now sees "src" as a package inside "BreadCounting"
command=/usr/bin/python3 -m uvicorn src.endpoint.Server:app --host 0.0.0.0 --port 8000
user=sunrise
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/uvicorn_err.log
stdout_logfile=/var/log/supervisor/uvicorn_out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5